# etc/sv/zookeeper/run
#!/bin/sh

echo " starting zookeeper"
exec /opt/zookeeper/bin/zkServer.sh start-foreground /etc/zookeeper/conf/zoo.cfg